# recipes-connectivity/wlan-udev-rules/files/wlan.rules
KERNEL=="wlan0", ACTION=="add" RUN+="/etc/udev/scripts/wlan.sh"
KERNEL=="wlan0", ACTION=="remove" RUN+="/etc/udev/scripts/wlan.sh" 
